# luci-app-mero/root/etc/uci-defaults/luci-mero-sync
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@mero[-1]
	add ucitrack mero
	set ucitrack.@mero[-1].init=mero
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
