# luci-app-mqos/root/etc/uci-defaults/mqos
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@mqos[-1]
	add ucitrack mqos
	set ucitrack.@mqos[-1].init=mqos
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0